# monitor-thumbnail.service
[Unit]
Description=Screenshot monitor stream

[Service]
ExecStart=/root/ffmpeg-screenshot-stream.sh
Type=oneshot

[Install]
WantedBy=timers.target
